# scripts/systemd/trainer-calendar-sync.service
[Unit]
Description=Sync trainer calendar from Google Calendar
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/home/kraftundfitness/website_K-F
ExecStart=/home/kraftundfitness/website_K-F/venv/bin/python /home/kraftundfitness/website_K-F/manage.py sync_trainer_calendar --force
