# scripts/systemd/trainer-calendar-sync.timer
[Unit]
Description=Run trainer calendar sync every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
Unit=trainer-calendar-sync.service

[Install]
WantedBy=timers.target

